# conf.d/stox.conf
[program:stox]
command=stox --no-sandbox
autorestart=true
