# conf.d/xfe.conf
[program:xfe]
command=xfe
autorestart=true
